# imageroot/systemd/user/twenty-app.service
#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

[Unit]
Description=Podman  twenty-app.service
BindsTo=twenty.service
After=twenty.service postgresql-app.service twenty-redis.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=%S/state/environment
EnvironmentFile=-%S/state/smarthost.env
WorkingDirectory=%S/state
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/twenty-app.pid %t/twenty-app.ctr-id
ExecStartPre=-runagent discover-smarthost
ExecStartPre=/bin/sh -c 'grep -q APP_SECRET environment || echo "APP_SECRET=`openssl rand -base64 32`" >> environment'
ExecStart=/usr/bin/podman run --conmon-pidfile %t/twenty-app.pid \
    --cidfile %t/twenty-app.ctr-id --cgroups=no-conmon \
    --pod-id-file %t/twenty.pod-id --replace -d --name  twenty-app \
    --volume twenty-app:/app/docker-data:Z \
    --env NODE_PORT=3000 \
    --env SERVER_URL=https://${TRAEFIK_HOST} \
    --env PG_DATABASE_URL="postgresql://postgres:Nethesis,1234@postgresql-app:5432/twenty" \
    --env REDIS_URL="redis://twenty-redis:6379" \
    --env APP_SECRET=${APP_SECRET} \
    --env NODE_ENV=production \
    --env LOG_LEVEL=info \
    --env STORAGE_TYPE=local \
    ${TWENTY_IMAGE}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/twenty-app.ctr-id -t 10
ExecReload=/usr/bin/podman kill -s HUP twenty-app
SyslogIdentifier=%u
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/twenty-app.ctr-id
PIDFile=%t/twenty-app.pid
Type=forking

[Install]
WantedBy=default.target
